# proxmigrate-backup.timer
[Unit]
Description=Ruleaza ProxMigrate Backup in fiecare zi

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true

[Install]
WantedBy=timers.target
